# log-generator.service
# EC2에서 systemd 서비스로 등록하는 파일
# 설치 방법:
#   sudo cp log-generator.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable log-generator
#   sudo systemctl start log-generator
#
# 모드 변경 방법:
#   /etc/log-generator/.env 의 LOG_MODE 값 수정 후
#   sudo systemctl restart log-generator

[Unit]
Description=Dummy Log Generator
After=network.target

[Service]
Type=simple
User=ec2-user
EnvironmentFile=/etc/log-generator/.env
ExecStartPre=/bin/mkdir -p /var/log/app
ExecStart=/usr/bin/python3 /opt/log-generator/log_generator.py \
    --mode ${LOG_MODE} \
    --interval ${LOG_INTERVAL} \
    --output ${LOG_OUTPUT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
